# matrix-sendmail.timer
[Unit]
Description=Timer to deliver email via matrix-sendmail.

[Timer]
OnUnitActiveSec=10m

